# templates/jenkins_agent.service
[Unit]
Description=Jenkins inbound agent
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/jenkins-agent
ExecStopPost=rm -rf /var/lib/jenkins/.ready

[Install]
WantedBy=multi-user.target
